# community/bspwm/profile.d/bspwm.sh
# Set our default path
PATH="/usr/local/sbin:/usr/local/bin:/usr/bin/core_perl:/usr/bin:$HOME/.config/bspwm/panel:$HOME/.bin"
export PANEL_FIFO="/tmp/panel-fifo"
export PATH
export BSPWM_SOCKET="/tmp/bspwm-socket"
export PANEL_HEIGHT=25
export GUI_EDITOR=/usr/bin/gedit
export BROWSER=/usr/bin/chromium
export TERMINAL=/usr/bin/sterminal
