# community/awesome/profile.d/awesome.sh

export GUI_EDITOR=/usr/bin/l3afpad
export BROWSER=/usr/bin/firefox
export TERMINAL=/usr/bin/sterminal
